# Reference/environments/c361768/modules/steng/files/etc/profile.d/jruby.sh
#!/bin/bash
export JRUBY_HOME="/opt/jruby"
PATH=${PATH}:${JRUBY_HOME}/bin
