# Reference/environments/c361768/modules/steng/files/backup/etc/profile.d/netbackup.sh
#!/bin/bash

# PATHS for NB7.1
# Excluded on purpose
#  /usr/openv/netbackup/bin/goodies
#  /usr/openv/netbackup/bin/support

PATH=$PATH:/usr/openv/netbackup/bin:/usr/openv/netbackup/bin/admincmd:/usr/openv/volmgr/bin
