# Reference/environments/c361768/modules/steng/files/backup/etc/profile.d/vrts.sh
#!/bin/bash

# PATHS for vrts

PATH=$PATH:/opt/VRTS/bin:/opt/VRTSvcs/bin:/usr/openv/netbackup/bin/admincmd/goodies